# assets/code/launchd/com.antcrate.intel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  antcrate :: launchd/com.antcrate.intel.plist — daily intel pull (retrieval
  only; classification stays human), macOS sibling of
  systemd/antcrate-intel.{service,timer}. Rendering notes: see backup plist.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antcrate.intel</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BIN__</string>
        <string>intel</string>
        <string>pull</string>
        <string>--quiet</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>__MIN__</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>ANTCRATE_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__LOG__/antcrate-intel.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG__/antcrate-intel.launchd.log</string>
</dict>
</plist>
