# assets/code/systemd/antcrate-backup.timer
[Unit]
Description=Daily AntCrate self-source backup timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
