# assets/code/systemd/antcrate-intel.service
[Unit]
Description=AntCrate daily Anthropic intel pull (retrieval only, no cognition)

[Service]
Type=oneshot
ExecStart=__BIN__ intel pull --quiet
Environment=ANTCRATE_LOG_LEVEL=info
